# framework/files/etc/systemd/system/rancherd.service.d/override.conf
[Unit]
ConditionPathExists=!/run/cos/live_mode
ConditionPathExists=!/run/cos/rescue_mode
