# vendor/github.com/coreos/coreos-cloudinit/units/user-cloudinit@.path
[Unit]
Description=Watch for a cloud-config at %f

[Path]
PathExists=%f
